# com.mcid.aio_pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mcid.aio_pipeline</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/thiagobruno/Desktop/GrandLine/NamiDash/painel_aio/aio_pipeline.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/thiagobruno/Desktop/GrandLine/NamiDash/painel_aio</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/thiagobruno/logs/aio_pipeline_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/thiagobruno/logs/aio_pipeline_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>TZ</key>
        <string>America/Sao_Paulo</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
